# scion/comn-gateway/envoy-comn.service
[Unit]
Description=Envoy COMN Gateway - LuciVerse "One Door"
Documentation=https://www.envoyproxy.io/docs
After=network-online.target
Wants=network-online.target
StartLimitBurst=3
StartLimitInterval=60s

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/envoy -c /etc/envoy/envoy.yaml --log-level info
Restart=on-failure
RestartSec=10

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Working directory
WorkingDirectory=/etc/envoy

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=envoy-comn

[Install]
WantedBy=multi-user.target
