# openeuler-embedded/meta-luciverse/recipes-security/lucitrust/files/lucitrust.service
[Unit]
Description=LuciTrust 7-Layer Trusted Computing System
Documentation=https://lucidigital.net/lucitrust
After=network.target tpm2-abrmd.service
Wants=tpm2-abrmd.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/python3 -m lucitrust.lucitrust_system --init --tier ${LUCIVERSE_TIER:-PAC}
ExecReload=/usr/bin/python3 -m lucitrust.lucitrust_system --init --tier ${LUCIVERSE_TIER:-PAC}

# Security hardening
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
ReadWritePaths=/var/lib/lucitrust

# Environment
Environment=PYTHONPATH=/usr/lib/lucitrust
Environment=LUCIVERSE_TIER=PAC
Environment=GENESIS_BOND_ID=GB-2025-0524-DRH-LCS-001

[Install]
WantedBy=multi-user.target
